# files/service/php-fpm/run
#!/bin/bash

/usr/bin/php-fpm -c /etc/php/fpm
